# .config/systemd/user/mastodon-redis.service
[Unit]
Description=mastodon redis instance

[Service]
Type=simple
ExecStart=/home/mastodon/run-redis.sh
