# .config/systemd/user/mastodon-databases.target
[Unit]
Description=Databases required by the Mastodon application
Requires=mastodon-postgres.service mastodon-redis.service
