# .config/systemd/user/mastodon-application.target
[Unit]
Description=Mastodon
Requires=mastodon-puma.service mastodon-sidekiq.target mastodon-streaming.service

[Install]
WantedBy=default.target
