# .config/systemd/user/mastodon-postgres.service
[Unit]
Description=mastodon postgres instance

[Service]
Type=simple
ExecStart=/home/mastodon/run-postgres.sh
